# src/etc/cron.d/cleanup-trash
# Delete mails older than ${C_DELETE_TRASH} days from all Trash folders
0 4 * * * find ${FRX_MAIL_DIR}/*/*/*/.Trash/cur/ -type f -mtime +${CLEANUP_TRASH} -exec rm {} \;
